# deploy/systemd-dropins/eddie-conversation-monitor.service.d/database.conf
[Service]
Environment="DATABASE_URL=postgresql://postgres:eddie_memory_2026@192.168.15.2:5432/postgres"
