# btc_trading_agent/systemd/crypto-exporter@SOL_USDT.service.d/env.conf
[Service]
Environment=METRICS_PORT=9095
Environment=COIN_CONFIG_FILE=config_SOL_USDT.json
