# grafana/exporters/trading-selfheal-exporter.service
[Unit]
Description=Trading Agent Runtime Self-Healing Exporter
Documentation=file:/apps/crypto-trader/trading/docs/TRADING_RUNTIME_SELFHEAL.md
After=network-online.target
Requires=network-online.target
StartLimitIntervalSec=120

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/apps/crypto-trader/trading/grafana/exporters
ExecStart=/usr/bin/python3 -u /apps/crypto-trader/trading/grafana/exporters/trading_selfheal_exporter.py \
  --port 9120 \
  --status-port 9121 \
  --config /apps/crypto-trader/trading/grafana/exporters/trading_selfheal_config.json \
  --interval 30

Restart=always
RestartSec=10
StartLimitBurst=5

Environment="DATABASE_URL=postgresql://postgres:eddie_memory_2026@192.168.15.2:5433/btc_trading"
Environment="OLLAMA_HOST=http://192.168.15.2:8512"
Environment="OLLAMA_MODEL=qwen2.5-coder:7b"
Environment="TRADING_HEAL_DATA_DIR=/var/lib/shared/trading-heal"

ReadWritePaths=/var/lib/shared/trading-heal /apps/crypto-trader/trading/btc_trading_agent/data /apps/crypto-trader/trading/btc_trading_agent/models /home/homelab
ProtectSystem=full
ProtectHome=no
NoNewPrivileges=false

StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-selfheal-exporter

[Install]
WantedBy=multi-user.target
