# systemd/ltfs-lto6.service.d/65-sg0-tape-access-scope.conf
[Service]
# sg0 usa fila/lock proprios. Isso evita acoplamento acidental com outros
# drives fisicos independentes.
Environment=TAPE_ACCESS_LOCKFILE=/run/lock/tape-access-sg0.lock
Environment=TAPE_ACCESS_QUEUE_DIR=/run/tape-queue-sg0
Environment=TAPE_ACCESS_HOLDER_FILE=/run/tape-queue-sg0/current
