# systemd/ltfs-lto6.service.d/80-selfheal-on-failure.conf
[Unit]
OnFailure=ltfs-lto6-selfheal-escalator.service
