# btc_trading_agent/systemd/crypto-agent@XRP_USDT.service.d/env.conf
[Service]
Environment=COIN_CONFIG_FILE=config_XRP_USDT.json
