# eddie-expurgo.service
[Unit]
Description=Eddie Gmail Expurgo Inteligente Service
After=network.target docker.service ollama.service

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/myClaude
Environment=PYTHONUNBUFFERED=1
Environment=OLLAMA_HOST=http://192.168.15.2:11434
Environment=WAHA_URL=http://localhost:3001
Environment=GMAIL_DATA_DIR=/home/homelab/myClaude/gmail_data
Environment=TELEGRAM_BOT_TOKEN=1105143633:AAEC1kmqDD_MDSpRFgEVHctwAfvfjVSp8B4
Environment=ADMIN_CHAT_ID=948686300
Environment=ADMIN_PHONE=5511999999999

ExecStart=/usr/bin/python3 /home/homelab/myClaude/gmail_expurgo_inteligente.py --execute --daemon --interval 24 --channels telegram whatsapp

Restart=always
RestartSec=60

# Logs
StandardOutput=append:/var/log/eddie-expurgo.log
StandardError=append:/var/log/eddie-expurgo-error.log

[Install]
WantedBy=multi-user.target
