# eddie-telegram-bot.service
[Unit]
Description=Eddie Coder Telegram Bot
After=network.target eddie-coder-api.service

[Service]
Type=simple
User=homelab
WorkingDirectory=/home/homelab/myClaude
Environment="PATH=/home/homelab/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="TELEGRAM_BOT_TOKEN=1105143633:AAEC1kmqDD_MDSpRFgEVHctwAfvfjVSp8B4"
Environment="OLLAMA_HOST=http://192.168.15.2:11434"
Environment="OLLAMA_MODEL=eddie-coder"
ExecStart=/usr/bin/python3 /home/homelab/myClaude/telegram_bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
