# location_integration/eddie-location.service
[Unit]
Description=Eddie Location Server - Integração de Localização
After=network.target

[Service]
Type=simple
User=eddie
WorkingDirectory=/home/eddie/myClaude/location_integration
Environment="PATH=/home/eddie/myClaude/location_integration/venv/bin"
# TELEGRAM_BOT_TOKEN must be provided at runtime from the vault or environment.
Environment="TELEGRAM_BOT_TOKEN="
Environment="TELEGRAM_CHAT_ID=948686300"
ExecStart=/home/eddie/myClaude/location_integration/venv/bin/python location_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
