# systemd/nvme-tape-drain.service.d/30-no-overlap.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/tape-session-guard --name nvme-tape-drain --busy-unit ltfs-cache-flush.service -- /usr/local/bin/nvme-tape-drain.sh
