# btc_trading_agent/systemd/crypto-agent@SOL_USDT.service.d/env.conf
[Service]
Environment=COIN_CONFIG_FILE=config_SOL_USDT.json
