# btc_trading_agent/systemd/crypto-exporter@ADA_USDT.service.d/env.conf
[Service]
Environment=METRICS_PORT=9097
Environment=COIN_CONFIG_FILE=config_ADA_USDT.json
