# systemd/nextcloud-tape-backup.timer.d/20-no-ltfs-require.conf
[Unit]
# O timer deve permanecer armado mesmo quando o mount LTFS ainda nao foi
# levantado. A dependencia real continua no service, nao no timer.
After=
Requires=
