# systemd/nvme-tape-drain.timer.d/20-no-ltfs-require.conf
[Unit]
# O timer deve continuar aguardando janela mesmo se o LTFS nao estiver montado.
# O service valida o mount real no momento da execucao.
After=
Requires=
