# monitoring/grafana/grafana-dashboard-sync.timer
[Unit]
Description=Periodically persist Grafana dashboard changes to provisioning files

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=60s
Persistent=true
Unit=grafana-dashboard-sync.service

[Install]
WantedBy=timers.target
