# tools/secrets_agent/secrets_agent.service
[Unit]
Description=Eddie Secrets Agent
After=network.target

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev
Environment=SECRETS_AGENT_API_KEY=please-set-a-strong-key
Environment=SECRETS_AGENT_DATA=/var/lib/eddie/secrets_agent
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python /home/edenilson/eddie-auto-dev/tools/secrets_agent/secrets_agent.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
