# tools/secrets_agent/secrets_agent.service
[Unit]
Description=Eddie Secrets Agent — gateway único para todos os secrets (v2.0 auto-unlock)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/agents_workspace/prod

# ── Configuração base ──
Environment=SECRETS_AGENT_API_KEY=please-set-a-strong-key
Environment=SECRETS_AGENT_DATA=/var/lib/eddie/secrets_agent
Environment=SECRETS_AGENT_PORT=8088
Environment=SECRETS_AGENT_PROM_PORT=8009

# ── Auto-unlock BW (sem solicitar senha) ──
# Opção 1: arquivo com master password (RECOMENDADO — mais seguro)
Environment=BW_PASSWORD_FILE=/var/lib/eddie/secrets_agent/.bw_master_password
# Opção 2: env var direta (menos seguro — visível em /proc)
# Environment=BW_MASTER_PASSWORD=sua-master-password
# Opção 3: API key login (para Bitwarden cloud)
# Environment=BW_CLIENTID=user.xxx
# Environment=BW_CLIENTSECRET=xxx
# Environment=BW_EMAIL=user@example.com

ExecStartPre=/bin/mkdir -p /var/lib/eddie/secrets_agent
ExecStart=/usr/bin/python3 /home/homelab/agents_workspace/prod/tools/secrets_agent/secrets_agent.py

# Always-on: reinicia em qualquer falha, sem limite
Restart=always
RestartSec=5
# Não habilitar watchdog systemd aqui sem implementar sd_notify(WATCHDOG=1),
# senão o serviço entra em loop de restart mesmo estando saudável.

# Health check removido do ExecStartPost (startup BW pode demorar)
# Use: curl http://127.0.0.1:8088/health para verificar manualmente

# Limites de recursos (evitar OOM)
MemoryMax=256M
MemoryHigh=192M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secrets-agent

[Install]
WantedBy=multi-user.target
