# tools/secrets_agent/secrets_agent.service
[Unit]
Description=Eddie Secrets Agent — gateway único para todos os secrets
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev
Environment=SECRETS_AGENT_API_KEY=please-set-a-strong-key
Environment=SECRETS_AGENT_DATA=/var/lib/eddie/secrets_agent
ExecStartPre=/bin/mkdir -p /var/lib/eddie/secrets_agent
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python /home/homelab/eddie-auto-dev/tools/secrets_agent/secrets_agent.py

# Always-on: reinicia em qualquer falha, sem limite
Restart=always
RestartSec=5
WatchdogSec=120

# Health check via curl (mata o serviço se /secrets não responder)
ExecStartPost=/bin/sh -c 'sleep 3 && curl -sf --connect-timeout 5 http://127.0.0.1:8088/secrets >/dev/null || exit 1'

# Limites de recursos (evitar OOM)
MemoryMax=256M
MemoryHigh=192M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secrets-agent

[Install]
WantedBy=multi-user.target
