# tools/x_agent/x-agent.service
[Unit]
Description=Eddie X Agent — X.com (Twitter) Integration
After=network-online.target secrets-agent.service
Wants=network-online.target

[Service]
Type=simple
User=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev
Environment=PYTHONUNBUFFERED=1
Environment=SECRETS_AGENT_URL=http://localhost:8088
Environment=SECRETS_AGENT_API_KEY=188bbf4c1b43ed1730005288f89ad2d0708c071eca142a2b335e026e95e8cee3
Environment=X_AGENT_PORT=8515
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python -m uvicorn tools.x_agent.x_agent:app --host 0.0.0.0 --port 8515
Restart=always
RestartSec=10
MemoryLimit=256M

[Install]
WantedBy=multi-user.target
