# btc_trading_agent/systemd/crypto-agent@ADA_USDT.service.d/env.conf
[Service]
Environment=COIN_CONFIG_FILE=config_ADA_USDT.json
