# systemd/crypto-agent@BTC_USDT_aggressive.service.d/zz-direct-ollama.conf
[Service]
# Roteamento via coordenador (11437) — GPU routing automático por modelo.
# Plano usa gemma3-fast (GPU1) para não competir com conservative/shadow no GPU0.
#   trading-analyst → GPU0 (11434, RTX 2060)
#   gemma3-fast:gpu1 → GPU1 (11435, GTX 1050)
Environment=OLLAMA_PLAN_HOST=http://192.168.15.2:11437
Environment=OLLAMA_PLAN_FALLBACK_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_PARAMS_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_PARAMS_FALLBACK_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_WINDOW_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_WINDOW_FALLBACK_HOST=http://192.168.15.2:11437
Environment=OLLAMA_PLAN_MODEL=gemma3-fast:gpu1
Environment=OLLAMA_PLAN_FALLBACK_MODEL=gemma3-fast:gpu1
