# btc_trading_agent/systemd/crypto-exporter@XRP_USDT.service.d/env.conf
[Service]
Environment=METRICS_PORT=9094
Environment=COIN_CONFIG_FILE=config_XRP_USDT.json
