# systemd/crypto-agent@.service.d/common.conf
[Unit]
After=docker.service eddie-postgres.service

[Service]
CPUAffinity=2-15
ExecStartPre=/usr/local/bin/wait_postgres.sh
# Centralizado: edite /etc/crypto-agent/models.env para trocar modelos
EnvironmentFile=/etc/crypto-agent/models.env
Environment=DATABASE_URL=postgresql://postgres:PASSWORD@localhost:5433/btc_trading
Environment=SECRETS_AGENT_API_KEY=<from_bitwarden>
Environment=TELEGRAM_BOT_TOKEN=<from_bitwarden>
Environment=ADMIN_CHAT_ID=<your_telegram_chat_id>
Environment=OLLAMA_PLAN_HOST=http://192.168.15.2:11434
Environment=OLLAMA_TRADE_PARAMS_HOST=http://192.168.15.2:11434
Environment=OLLAMA_TRADE_PARAMS_FALLBACK_HOST=http://192.168.15.2:11435
Environment=OLLAMA_TRADE_PARAMS_MODE=apply
Environment=OLLAMA_TRADE_WINDOW_HOST=http://192.168.15.2:11434
Environment=OLLAMA_TRADE_WINDOW_FALLBACK_HOST=http://192.168.15.2:11435
Restart=always
RestartSec=30
