# systemd/crypto-agent@.service.d/deps.conf
[Unit]
After=docker.service eddie-postgres.service
Wants=eddie-postgres.service

[Service]
ExecStartPre=/usr/local/bin/wait_postgres.sh
