# systemd/crypto-agent@.service.d/ollama-timeout.conf
[Service]
# Timeouts baseados em velocidade real dos modelos:
# trading-analyst ~6 tok/s × 1024 tokens = 170s → 200s com margem
# qwen3-fast ~58 tok/s mas fila coordinator até 60s → 120s cobre
Environment="OLLAMA_TRADE_WINDOW_TIMEOUT_SEC=120"
Environment="OLLAMA_TRADE_WINDOW_FALLBACK_TIMEOUT_SEC=120"
Environment="OLLAMA_PLAN_TIMEOUT_SEC=200"
Environment="OLLAMA_TRADE_PARAMS_TIMEOUT_SEC=120"
Environment="OLLAMA_TRADE_PARAMS_FALLBACK_TIMEOUT_SEC=90"
