# systemd/crypto-agent@.service.d/zz-ai-throttle.conf
[Service]
# Contencao de chamadas Ollama apos upgrade de GPU.
# RSS pode publicar varias noticias em rajada e disparar plano+controls+window
# nos 3 perfis ao mesmo tempo. Mantem trading ativo e reduz uso de LLM.
Environment=OLLAMA_RSS_TRIGGERS_ENABLED=false
Environment=OLLAMA_AI_PLAN_MIN_INTERVAL_SEC=900
Environment=OLLAMA_TRADE_PARAMS_MIN_INTERVAL_SEC=1800
Environment=OLLAMA_TRADE_WINDOW_MIN_INTERVAL_SEC=300
Environment=OLLAMA_TRADE_WINDOW_MIN_INTERVAL_AGGRESSIVE_SEC=240
Environment=OLLAMA_TRADE_WINDOW_MIN_INTERVAL_CONSERVATIVE_SEC=420
Environment=OLLAMA_TRADE_WINDOW_TTL_SEC=300
Environment=OLLAMA_TRADE_WINDOW_TTL_AGGRESSIVE_SEC=240
Environment=OLLAMA_TRADE_WINDOW_TTL_CONSERVATIVE_SEC=420
