# btc_trading_agent/systemd/crypto-exporter@ETH_USDT.service.d/env.conf
[Service]
Environment=METRICS_PORT=9093
Environment=COIN_CONFIG_FILE=config_ETH_USDT.json
