# systemd/wg-quick@protonvpn.service.d/restore-iprules.conf
[Service]
ExecStartPost=/bin/bash -c 'sleep 90 && /usr/local/bin/protonvpn-routing-watchdog.sh --ensure || true &'
