# systemd/akash-sweep.service.d/secrets.conf
[Service]
# Chave de API do Secrets Agent — obter via:
# ssh homelab@192.168.15.2 "sudo systemctl cat btc-trading-agent | grep SECRETS_AGENT_API_KEY"
# ou via Bitwarden: authentik/shared/secrets_agent_api_key
Environment=SECRETS_AGENT_API_KEY=CHANGEME
Environment=SECRETS_AGENT_URL=http://192.168.15.2:8088
