# systemd/lto6-drain-backups.service.d/50-tape-gate.conf
# Serializa acesso à fita: drain adquire lock EXCLUSIVO antes de rodar.
# Se tape-read-wrap (leitura direta da fita) estiver ativo, AGUARDA ao invés de falhar.
# Substitui o ExecStartPre simplista que apenas falhava se flush estava ativo.
[Service]
ExecStartPre=
ExecStartPre=/bin/bash -c \
  "while systemctl is-active --quiet ltfs-cache-flush.service; do \
     echo 'ltfs-cache-flush ativo, aguardando 30s...'; sleep 30; \
   done"
ExecStart=
ExecStart=/usr/local/sbin/tape-exclusive-wrap -- /usr/local/sbin/lto6-drain-backups
