# systemd/ollama-gpu1.service.d/memory-limit.conf
[Service]
# gemma3-fast:gpu1 usa ~1.6GB de VRAM, mas o processo Ollama + runner pode
# passar de 3GB de RAM durante load/reload. O limite antigo empurrava swap e
# deixava a GPU1 instavel sob concorrencia.
MemoryHigh=6G
MemoryMax=8G
