# systemd/ollama-gpu1.service.d/zzzz-warmup-curl.conf
[Service]
# Warmup nao interativo e limitado por timeout para a GPU1.
ExecStartPost=
ExecStartPost=/usr/bin/python3 /apps/crypto-trader/tools/ollama_warmup.py --host http://127.0.0.1:11435 --default-model gemma3-fast:gpu1 --timeout 90 --delay 8
