# systemd/ollama-gpu-coordinator.service.d/zz-dual-gpu-routing.conf
[Service]
# Roteamento pós-upgrade:
# - GPU0/11434: RTX 3060 12GB para trading-analyst e modelos pesados.
# - GPU1/11435: RTX 2060 SUPER 8GB para gemma3-fast:gpu1 e janelas rápidas.
Environment=OLLAMA_GPU0_HOST=http://192.168.15.2:11434
Environment=OLLAMA_GPU1_HOST=http://192.168.15.2:11435
# O busy check atual usa /api/ps, que mostra modelos keep_alive como "ocupados".
# Ate a Fase 2 trocar isso por semaforo real, nao aguardar artificialmente.
Environment=GPU_COORD_BUSY_WAIT_SEC=0
