# systemd/banking-metrics-exporter.service
[Unit]
Description=Banking Agent Metrics Exporter (Prometheus)
After=network.target eddie-postgres.service specialized-agents-api.service
Wants=eddie-postgres.service

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev
EnvironmentFile=-/etc/default/banking-metrics-exporter
Environment="BANKING_EXPORTER_PORT=9102"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python -m specialized_agents.banking_metrics_exporter
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
