# systemd/candle-collector.service
[Unit]
Description=Coleta periódica de candles KuCoin → PostgreSQL
After=network.target postgresql.service

[Service]
Type=oneshot
User=homelab
WorkingDirectory=/home/homelab/myClaude/btc_trading_agent
Environment=PYTHONUNBUFFERED=1
Environment=DATABASE_URL=postgresql://postgres:eddie_memory_2026@localhost:5433/btc_trading
Environment=PYTHONPATH=/home/homelab/eddie-auto-dev

ExecStart=/usr/bin/python3 /home/homelab/eddie-auto-dev/scripts/candle_collector.py --hours 2 --symbols BTC-USDT ETH-USDT SOL-USDT XRP-USDT DOGE-USDT ADA-USDT

StandardOutput=journal
StandardError=journal
SyslogIdentifier=candle-collector

TimeoutStartSec=120
