# systemd/crypto-agent@.service
[Unit]
Description=Crypto Trading Agent - %I
After=network.target postgresql.service secrets_agent.service
Wants=network-online.target secrets_agent.service
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=simple
User=trading-svc
Group=trading-svc
WorkingDirectory=/apps/crypto-trader/trading/btc_trading_agent
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=COIN_CONFIG_FILE=config_%I.json
Environment=HOME=/apps/crypto-trader

# Porta de métricas/API configurável via override
# Cada instância DEVE ter porta única (ver deploy_profiles.sh)
EnvironmentFile=-/apps/crypto-trader/envfiles/%I.env

ExecStart=/usr/bin/python3 /apps/crypto-trader/trading/btc_trading_agent/trading_agent.py --daemon --live --config config_%I.json
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-agent-%I

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/apps/crypto-trader/trading/btc_trading_agent /apps/crypto-trader/data

[Install]
WantedBy=multi-user.target
