# systemd/crypto-agent@.service
[Unit]
Description=Crypto Trading Agent - %I
After=network.target postgresql.service secrets_agent.service
Wants=network-online.target secrets_agent.service

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/myClaude/btc_trading_agent
Environment=PATH=/home/homelab/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=COIN_CONFIG_FILE=config_%I.json

# Porta de métricas/API configurável via override
# Cada instância DEVE ter porta única (ver deploy_profiles.sh)
EnvironmentFile=-/home/homelab/myClaude/btc_trading_agent/envfiles/%I.env

ExecStart=/usr/bin/python3 /home/homelab/myClaude/btc_trading_agent/trading_agent.py --daemon --live --config config_%I.json
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-agent-%I

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/homelab/myClaude/btc_trading_agent

[Install]
WantedBy=multi-user.target
