# systemd/disk-clean.timer
[Unit]
Description=Weekly disk cleanup timer

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
