# systemd/eddie-telegram-bot.service
[Unit]
Description=Eddie Telegram Bot — chatbot com trading report, home assistant, Gmail e agentes especializados
After=network.target specialized-agents-api.service

[Service]
Type=simple
WorkingDirectory=/home/homelab/eddie-auto-dev
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python3 /home/homelab/myClaude/telegram_bot.py
Restart=on-failure
RestartSec=10
User=homelab
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/eddie-common

[Install]
WantedBy=multi-user.target
