# systemd/homelab-tape-log-drain-nextcloud.service
[Unit]
Description=Drain buffered logs stamped for the Nextcloud tape pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=ROUTE_NAME=nextcloud
Environment=ROUTE_TARGET_ROOT=/mnt/raid1/lto6-cache/logs
ExecStart=/usr/local/sbin/tape-log-spool-drain
User=root
Group=root
