# systemd/homelab-tape-log-drain-nextcloud.timer
[Unit]
Description=Drain buffered logs stamped for the Nextcloud tape pipeline regularly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10min
Persistent=true
Unit=homelab-tape-log-drain-nextcloud.service

[Install]
WantedBy=timers.target
