# systemd/kucoin-usdt-rebalance.timer
[Unit]
Description=KuCoin USDT rebalance — a cada hora

[Timer]
OnCalendar=hourly
Persistent=true
Unit=kucoin-usdt-rebalance.service

[Install]
WantedBy=timers.target
