# systemd/ltfs-lto6-sg1.service
[Unit]
Description=LTFS mount for HP LTO-6 on sg1
After=local-fs.target network.target
Requires=local-fs.target
Conflicts=tape-safe-eject.service ltfs-deep-recovery.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/ltfs-lto6-sg1
ExecStart=/usr/bin/python3 /usr/local/tools/ltfs_recovery.py --orchestrated-mount
ExecStop=/usr/bin/python3 /usr/local/tools/ltfs_recovery.py --orchestrated-stop
ExecStartPost=/bin/mkdir -p /run/ltfs-export/lto6-sg1
ExecStartPost=/bin/mount --bind /mnt/tape/lto6-sg1 /run/ltfs-export/lto6-sg1
ExecStopPost=-/bin/umount -l /run/ltfs-export/lto6-sg1
ExecStopPost=-/bin/rmdir /run/ltfs-export/lto6-sg1
TimeoutStartSec=480
TimeoutStopSec=120
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
