# systemd/lto-logical-mount-refresh.service
[Unit]
Description=Refresh the logical single-mount view for LTO cache and tape targets
After=network-online.target remote-fs.target
Wants=network-online.target
RequiresMountsFor=/mnt/raid1
ConditionPathExists=/usr/bin/mergerfs

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=-/etc/default/lto-logical-mount
ExecStart=/usr/local/bin/lto-logical-mount refresh
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lto-logical-mount
