# systemd/lto6-drain-backups.service
[Unit]
Description=Drain complete RPA4ALL backup snapshots to LTO-6 tape
Documentation=https://github.com/rpa4all/eddie-auto-dev
After=network-online.target mnt-lto6\x2dsmb\x2dproof.mount
Wants=mnt-lto6\x2dsmb\x2dproof.mount
ConditionPathIsMountPoint=/mnt/lto6-smb-proof

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/lto6-drain-backups
StandardOutput=journal
StandardError=journal
TimeoutStartSec=7200
# Prevent competing with LTFS cache flush
ExecStartPre=/bin/bash -c "systemctl is-active ltfs-cache-flush.service && { echo 'ltfs-cache-flush running, waiting'; exit 1; } || true"
