# systemd/marketing-x-posts.service
[Unit]
Description=Post orgânico no X/Twitter — Marketing RPA4ALL
After=network.target postgresql.service

[Service]
Type=oneshot
User=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/homelab/eddie-auto-dev
Environment=DATABASE_URL=postgresql://postgres:postgres@192.168.15.2:5433/shared
Environment=X_AGENT_URL=http://localhost:8515

ExecStart=/usr/bin/python3 /home/homelab/eddie-auto-dev/marketing/x_post_scheduler.py

StandardOutput=journal
StandardError=journal
