# systemd/marketing-x-posts.timer
[Unit]
Description=Timer para posts orgânicos no X (Ter/Qui/Sáb às 10:30)

[Timer]
OnCalendar=Tue *-*-* 10:30:00
OnCalendar=Thu *-*-* 10:30:00
OnCalendar=Sat *-*-* 10:30:00
Persistent=true

[Install]
WantedBy=timers.target
