# systemd/ollama-warmup.service
[Unit]
Description=Ollama GPU Warmup — mantém eddie-sentiment carregado na GPU0
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=homelab
WorkingDirectory=/home/homelab/ollama-warmup
ExecStart=/usr/bin/python3 /home/homelab/ollama-warmup/ollama_warmup.py --verbose
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
