# systemd/ollama-finetune.service
[Unit]
Description=Ollama Fine-Tuning Batch — QLoRA eddie-sentiment
After=network-online.target
Wants=network-online.target
Conflicts=ollama-warmup.service

[Service]
Type=oneshot
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev

ExecStart=/opt/finetune-env/bin/python3 /home/homelab/eddie-auto-dev/scripts/ollama_finetune_batch.py

# Fine-tuning pode levar até 2h
TimeoutStartSec=7200

Environment=OLLAMA_HOST=http://192.168.15.2:11434
Environment=OLLAMA_HOST_GPU1=http://192.168.15.2:11435
Environment=DATABASE_URL=postgresql://postgres:eddie_memory_2026@192.168.15.2:5433/btc_trading
Environment=CUDA_VISIBLE_DEVICES=0
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama-finetune

[Install]
WantedBy=multi-user.target
