# systemd/ollama-warmup.service
[Unit]
Description=Ollama GPU Warmup — mantém modelos configurados carregados nas GPUs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=edenilson
WorkingDirectory=/home/edenilson/eddie-auto-dev
Environment=OLLAMA_WARM_MODEL_GPU0=trading-analyst:latest
Environment=OLLAMA_WARM_MODEL_GPU1=qwen3:0.6b
ExecStart=/home/edenilson/eddie-auto-dev/.venv/bin/python3 /home/edenilson/eddie-auto-dev/scripts/ollama_warmup.py --verbose
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
