# systemd/pihole-ipv6-dns-fix.service
[Unit]
Description=Pi-hole DNS Intercept — força todo DNS (IPv4+IPv6) para Pi-hole
Documentation=https://github.com/eddiejdi/eddie-auto-dev
# Garante execução após rede E Docker estarem prontos
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
# Reinicia se o Docker reiniciar (ex: após reboot do container pihole)
PartOf=docker.service
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/pihole-ipv6-dns-fix.sh
# Tenta novamente se falhar (ex: container ainda subindo)
Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
