# systemd/pihole-ipv6-dns-fix.service
[Unit]
Description=Pi-hole IPv6 DNS Interception Fix
Documentation=https://github.com/eddiejdi/eddie-auto-dev
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole-ipv6-dns-fix.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
