# systemd/rag-reindex.timer
[Unit]
Description=Daily RAG reindex timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
