# systemd/rss-sentiment-exporter.service
[Unit]
Description=RSS Sentiment Exporter — Crypto News Sentiment (Prometheus)
After=network.target eddie-postgres.service
Wants=eddie-postgres.service

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev
EnvironmentFile=-/etc/default/rss-sentiment-exporter
Environment="DATABASE_URL=postgresql://postgres:shared_memory_2026@192.168.15.2:5433/btc_trading"
Environment="OLLAMA_HOST_GPU1=http://192.168.15.2:11435"
Environment="OLLAMA_HOST=http://192.168.15.2:11434"
Environment="OLLAMA_SENTIMENT_MODEL=qwen2.5-coder:7b"
Environment="RSS_FETCH_INTERVAL=300"
Environment="RSS_SENTIMENT_WINDOW=4"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python grafana/exporters/rss_sentiment_exporter.py --port 9122
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
