# systemd/tape-component-quality-exporter.service
[Unit]
Description=Tape Component Quality Exporter
After=network.target
Wants=network.target
Conflicts=ltfs-deep-recovery.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/tape-component-quality-exporter
Environment="PYTHONUNBUFFERED=1"
Environment="TAPE_COMPONENT_QUALITY_PORT=9124"
ExecStart=/usr/bin/python3 /usr/local/tools/tape_component_quality_agent.py --exporter --port ${TAPE_COMPONENT_QUALITY_PORT} --interval 30
Restart=always
RestartSec=30
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tape-component-quality-exporter

# Precisa de acesso local a /sys/class/fc_host e /dev/sg*
NoNewPrivileges=no
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
