# systemd/trading-analyst-weekly-retrain.service
[Unit]
Description=Retreino semanal do trading-analyst-candidate (dataset->treino->shadow->relatorio)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=homelab
Group=homelab
# Janela ampla: o ciclo (dataset+treino+shadow) leva ~30-60min.
TimeoutStartSec=7200
Environment=PYTHONUNBUFFERED=1
# Fonte canônica de DATABASE_URL + TELEGRAM_BOT_TOKEN (garante o envio do relatório).
EnvironmentFile=/etc/default/eddie-common
ExecStart=/home/homelab/finetune/trading_analyst_weekly_retrain.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
