# systemd/tuya-token-selfheal.timer
[Unit]
Description=Roda o self-heal do token Tuya a cada 5 minutos
After=network-online.target

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=tuya-token-selfheal.service

[Install]
WantedBy=timers.target
