# systemd/workstation-win11l-rdp@.service
[Unit]
Description=On-demand proxy for Windows 11 Light Workstation RDP
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
ExecStartPre=/usr/local/bin/workstation-win11l-run.sh
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do nc -z 127.0.0.1 13391 && exit 0; sleep 2; done; exit 1'
ExecStart=/usr/bin/nc 127.0.0.1 13391
StandardInput=socket
StandardOutput=socket
PrivateTmp=true