# deploy/vpn/nordvpn-routing-watchdog.service
[Unit]
Description=NordVPN Routing Watchdog Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nordvpn-routing-watchdog.sh --ensure
Environment=HOME=/root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nordvpn-watchdog

[Install]
WantedBy=multi-user.target
