# deploy/vpn/protonvpn-routing-watchdog.service
[Unit]
Description=ProtonVPN Routing Watchdog Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/local/bin/protonvpn-routing-watchdog.sh --ensure
Environment=HOME=/root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=protonvpn-watchdog

[Install]
WantedBy=multi-user.target
