# deploy/vpn/rpa4all-ddns-server.timer
[Unit]
Description=Timer DDNS Cloudflare — atualiza vpn.rpa4all.com a cada 5min

[Timer]
OnBootSec=30
OnUnitActiveSec=5min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
