# deploy/vpn-deb/rpa4all-vpn/usr/share/rpa4all-vpn/rpa4all-vpn-ddns.service
[Unit]
Description=Atualizar endpoint DDNS da VPN RPA4All

[Service]
Type=oneshot
ExecStart=/usr/bin/rpa4all-vpn-update-endpoint
