# deploy/vpn-deb/rpa4all-vpn/usr/share/rpa4all-vpn/rpa4all-vpn-ddns.timer
[Unit]
Description=Timer DDNS para VPN RPA4All (a cada 30min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
