# btc_trading_agent/btc-daily-report.service
[Unit]
Description=Bitcoin Trading Agent - Daily Report
After=network.target waha.service

[Service]
Type=oneshot
User=homelab
WorkingDirectory=/home/homelab/myClaude/btc_trading_agent
ExecStart=/usr/bin/python3 /home/homelab/myClaude/btc_trading_agent/daily_report.py
StandardOutput=append:/home/homelab/myClaude/btc_trading_agent/logs/daily_report.log
StandardError=append:/home/homelab/myClaude/btc_trading_agent/logs/daily_report.log

[Install]
WantedBy=multi-user.target
