# systemd/ltfs-cache-flush.service.d/70-rearm-timer-on-exit.conf
# Rearm ltfs-cache-flush.timer apos cada execucao (sucesso ou falha).
# Garante que o timer nao morra silenciosamente se o servico for ativado
# diretamente (ex: por idle-unmount ou escalator), perdendo a proxima janela.
[Service]
ExecStopPost=
ExecStopPost=/bin/systemctl --no-block restart ltfs-cache-flush.timer
