# tools/systemd/eddie-weather-agent.service
[Unit]
Description=Eddie Weather Monitoring Agent — coleta meteorológica a cada 15 min
After=network-online.target postgresql.service eddie-postgres.service
Wants=network-online.target

[Service]
Type=simple
User=edenilson
Group=edenilson
WorkingDirectory=/home/edenilson/eddie-auto-dev
Environment=DATABASE_URL=postgresql://postgres:eddie_memory_2026@localhost:55432/postgres
Environment=WEATHER_LATITUDE=-23.5505
Environment=WEATHER_LONGITUDE=-46.6333
Environment=WEATHER_LOCATION=São Paulo, BR
Environment=WEATHER_TIMEZONE=America/Sao_Paulo
Environment=WEATHER_INTERVAL=900
ExecStart=/home/edenilson/eddie-auto-dev/.venv/bin/python tools/weather_agent.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
