# btc_trading_agent/systemd/crypto-agent@ETH_USDT.service.d/env.conf
[Service]
Environment=COIN_CONFIG_FILE=config_ETH_USDT.json
