# systemd/ollama.service.d/zz-gpu0-visible-device.conf
[Service]
# Isola a RTX 3060 como GPU0. Usar o indice evita fallback para CPU observado
# no Ollama 0.17.6 quando CUDA_VISIBLE_DEVICES recebeu o UUID pos-upgrade.
Environment=CUDA_VISIBLE_DEVICES=0
