# systemd/ollama.service.d/zzzz-warmup-curl.conf
[Service]
# Warmup principal desativado temporariamente: o trading-analyst carregou em
# VRAM, mas a primeira geracao ficou presa em "llm server loading model".
# Mantem o servico subindo rapido; a carga passa a ser sob demanda.
ExecStartPost=
ExecStartPost=/bin/true
