# btc_trading_agent/systemd/crypto-api@XRP_USDT.service.d/env.conf
[Service]
Environment=BTC_ENGINE_API_PORT=8513
Environment=COIN_CONFIG_FILE=config_XRP_USDT.json
