# btc_trading_agent/systemd/crypto-api@ADA_USDT.service.d/env.conf
[Service]
Environment=BTC_ENGINE_API_PORT=8516
Environment=COIN_CONFIG_FILE=config_ADA_USDT.json
