# systemd/crypto-agent@BTC_USDT_conservative.service.d/ollama-routing.conf
[Service]
# Roteamento via coordenador para serializar com os demais agentes.
# Coordenador (11437) gerencia GPU0/GPU1 automaticamente.
Environment=OLLAMA_TRADE_WINDOW_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_WINDOW_FALLBACK_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_PARAMS_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_PARAMS_FALLBACK_HOST=http://192.168.15.2:11437
