# systemd/crypto-agent@BTC_USDT_conservative.service.d/zz-direct-ollama.conf
[Service]
# Roteamento via coordenador (11437) — GPU routing automático por modelo:
#   trading-analyst → GPU0 (11434, RTX 2060)
#   gemma3-fast:gpu1 → GPU1 (11435, GTX 1050)
# Coordenador serializa acessos e evita 503 em cascata.
Environment=OLLAMA_PLAN_HOST=http://192.168.15.2:11437
Environment=OLLAMA_PLAN_FALLBACK_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_PARAMS_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_PARAMS_FALLBACK_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_WINDOW_HOST=http://192.168.15.2:11437
Environment=OLLAMA_TRADE_WINDOW_FALLBACK_HOST=http://192.168.15.2:11437
