# tools/selfheal/mnt-lto6-smb-proof.mount
[Unit]
Description=CIFS proof mount — LTO6 NAS reachability gate
Documentation=https://github.com/rpa4all/eddie-auto-dev
After=network-online.target
Wants=network-online.target

[Mount]
What=//192.168.15.4/LTO6_CACHE
Where=/mnt/lto6-smb-proof
Type=cifs
Options=credentials=/root/.smb-lto6-cache-credentials,vers=3.0,iocharset=utf8,uid=1000,gid=1000,file_mode=0664,dir_mode=0775,soft,nounix,_netdev
TimeoutSec=30

[Install]
WantedBy=multi-user.target
