# btc_trading_agent/systemd/crypto-api@SOL_USDT.service.d/env.conf
[Service]
Environment=BTC_ENGINE_API_PORT=8514
Environment=COIN_CONFIG_FILE=config_SOL_USDT.json
