# btc_trading_agent/systemd/crypto-api@ETH_USDT.service.d/env.conf
[Service]
Environment=BTC_ENGINE_API_PORT=8512
Environment=COIN_CONFIG_FILE=config_ETH_USDT.json
