# tools/homelab_recovery/operations-agent.service
[Unit]
Description=Operations Agent - RCA remediation handler
After=network.target agent-api.service

[Service]
Type=simple
WorkingDirectory=/home/homelab/eddie-auto-dev
Environment="PYTHONPATH=/home/homelab/eddie-auto-dev"
Environment="AGENT_API_URL=http://127.0.0.1:8888"
Environment="ALLOW_AGENT_API=1"
Environment="DATABASE_URL=postgresql://postgres:eddie_memory_2026@localhost:5432/postgres"
Environment="OPS_AGENT_POLL=10"
Environment="AUTONOMOUS_MODE=0"
ExecStart=/home/homelab/eddie-auto-dev/.venv/bin/python3 -u /home/homelab/eddie-auto-dev/tools/operations_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
