# specialized_agents/specialized-agents-api.service
[Unit]
Description=Specialized Agents API (On-Demand)
After=network.target docker.service

[Service]
Type=simple
User=homelab
WorkingDirectory=/home/homelab/myClaude/specialized_agents
Environment=PATH=/home/homelab/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/homelab/myClaude

# Configurações On-Demand via ambiente
Environment=ONDEMAND_ENABLED=true
Environment=ONDEMAND_AGENT_TIMEOUT=300
Environment=ONDEMAND_DOCKER_TIMEOUT=600
Environment=ONDEMAND_RAG_TIMEOUT=300
Environment=ONDEMAND_GITHUB_TIMEOUT=180
Environment=ONDEMAND_CLEANUP_INTERVAL=60

# Carregar variáveis do .env se existir
EnvironmentFile=-/home/homelab/myClaude/.env

ExecStart=/home/homelab/.local/bin/uvicorn api_ondemand:app --host 0.0.0.0 --port 8503
Restart=always
RestartSec=10

# Limites de recursos (leve pois é on-demand)
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
