# btc_trading_agent/systemd/crypto-api@DOGE_USDT.service.d/env.conf
[Service]
Environment=BTC_ENGINE_API_PORT=8515
Environment=COIN_CONFIG_FILE=config_DOGE_USDT.json
