# tools/tunnels/cloudflared-named@.service
[Unit]
Description=Cloudflared named tunnel service for %i
After=network.target

[Service]
Type=simple
# We keep this unit simple — credentials and config are expected under /etc/cloudflared
ExecStart=/usr/bin/cloudflared tunnel run %i --config /etc/cloudflared/config.yml
Restart=on-failure
RestartSec=5
User=root
WorkingDirectory=/var/lib/cloudflared

[Install]
WantedBy=multi-user.target
