# systemd/coordinator-agent.service.d/coordinator-version.conf
# Feature flag para migração LangGraph (Fase 4)
# Alterar para v2 para ativar coordinator_v2.py
# Rollback imediato: mudar para v1 + systemctl restart coordinator-agent
[Service]
Environment="COORDINATOR_VERSION=v1"
EnvironmentFile=-/etc/default/eddie-common
