# systemd/nextcloud-tape-backup.service.d/30-no-overlap.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/tape-session-guard --name nextcloud-tape-backup --busy-unit ltfs-cache-flush.service -- /usr/local/sbin/nextcloud-tape-backup.sh
