# systemd/nextcloud-tape-backup.service.d/35-sg0-tape-access-scope.conf
[Service]
Environment=TAPE_ACCESS_LOCKFILE=/run/lock/tape-access-sg0.lock
Environment=TAPE_ACCESS_QUEUE_DIR=/run/tape-queue-sg0
Environment=TAPE_ACCESS_HOLDER_FILE=/run/tape-queue-sg0/current
